# docker/django/supervisor-app.conf
[supervisord]
[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /srv/code/configs/uwsgi.ini
